# com.familyfinance.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.familyfinance.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-c</string>
        <string>import sys; sys.path.insert(0, '/Users/suren/Library/Python/3.9/lib/python/site-packages'); import uvicorn; uvicorn.run('main:app', host='0.0.0.0', port=8000)</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/suren/family-finance/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/suren/family-finance/logs/app.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/suren/family-finance/logs/app.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
